# scripts/crm-importacion-runner.service
[Unit]
# LEGACY: este service era del modelo host-native. Hoy el runner corre como
# container Docker (fidcore-crm-importacion-runner). Se mantiene solo como
# fallback si alguna vez se revierte a host-native. NO instalar por default.
# Antes de usarlo, editar User y las 2 rutas absolutas con los valores reales
# de la instalación (típico: User=<usuario>, path=/opt/crm-fidcore).
Description=CRM FidCore - Importacion jobs runner loop (cada 30s)
After=network.target

[Service]
Type=simple
User=CHANGEME_USUARIO
WorkingDirectory=/opt/crm-fidcore
ExecStart=/bin/bash /opt/crm-fidcore/scripts/importacion-runner-loop.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
